# scripts/security/rotation-scheduling/credential-rotation.timer
[Unit]
Description=Credential Rotation Timer (Monthly)
Requires=credential-rotation.service

[Timer]
# Run on the first day of every month at 2 AM
OnCalendar=monthly
# Add randomization to avoid all systems rotating at once
RandomizedDelaySec=3600

# Run immediately if missed (e.g., system was down)
Persistent=true

# Ensure timer is accurate
AccuracySec=1h

[Install]
WantedBy=timers.target